# dashboard-web.service
[Unit]
Description=Arthas Trading Dashboard — Next.js web frontend
After=network.target dashboard-api.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/trading/web
Environment=PORT=3000
Environment=NEXT_PUBLIC_API_URL=http://187.77.5.109:8000
ExecStart=/usr/bin/node /opt/trading/web/node_modules/.bin/next start -p 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
