# stocks-agent.service
[Unit]
Description=Arthas Stocks Agent — NYSE/NASDAQ Momentum Trading
After=network.target postgresql.service telegram-bot.service
Wants=postgresql.service

[Service]
Type=simple
WorkingDirectory=/opt/trading
ExecStart=/opt/trading/venv/bin/python3 /opt/trading/scripts/run_stocks.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

# Limitar CPU al 50% para no competir con el trading agent de crypto
CPUQuota=50%

# Reiniciar en cualquier fallo
StartLimitIntervalSec=60
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
